# deploy/uniclaw.service
[Unit]
Description=UniClaw AI Assistant
After=cloud-sql-proxy.service docker.service
Requires=cloud-sql-proxy.service

[Service]
Type=simple
ExecStartPre=/usr/bin/docker pull us-central1-docker.pkg.dev/uniclaw-prod/uniclaw/agent:latest
ExecStart=/usr/bin/docker run --rm \
  --name uniclaw \
  --env-file /opt/uniclaw/.env \
  --network=host \
  us-central1-docker.pkg.dev/uniclaw-prod/uniclaw/agent:latest \
  --no-onboard
ExecStop=/usr/bin/docker stop uniclaw
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
